# tools/overlay_ubuntu22_04/etc/rc.local
#!/bin/bash
[ -f '/etc/rc.local.d/rcS' ] && bash /etc/rc.local.d/rcS







exit 0
