# tools/overlay_ubuntu22_04_LITE/etc/rc.local
#!/bin/bash
# /usr/local/m5stack/bin/S03espc6wifi start

[ -f '/etc/rc.local.d/rcS' ] && bash /etc/rc.local.d/rcS






exit 0
